# scripts/paper-pdf-summary/deploy/paper-pdf-summary-telegram.service
[Unit]
Description=Paper PDF Summary Telegram Bot
After=network.target

[Service]
Type=simple
User=xulei
WorkingDirectory=/opt/lis-rss-daily/scripts/paper-pdf-summary
ExecStart=/home/xulei/.pyenvs/env_camoufox/bin/python -m telegram_bot.main
Restart=always
RestartSec=5
Environment="PATH=/home/xulei/.pyenvs/env_camoufox/bin"
EnvironmentFile=/opt/lis-rss-daily/scripts/paper-pdf-summary/.env

[Install]
WantedBy=multi-user.target